# example/etc/udev/rules.d/99-hems.rules
ACTION=="add",ENV{ID_BUS}=="usb",ENV{DEVNAME}=="/dev/ttyACM0",RUN+="/bin/systemctl --no-block start hems.service"
